# runner/docker-gc.cron
# See the meaning of the provided vars at https://github.com/spotify/docker-gc#manual-usage
# m h dom mon dow command
0 * * * * DOCKER=balena-engine REMOVE_VOLUMES=1 FORCE_IMAGE_REMOVAL=1 FORCE_CONTAINER_REMOVAL=1 /usr/bin/docker-gc > /var/log/docker-gc.log
