# runner/docker-gc.cron
# See the meaning of the provided vars at https://github.com/spotify/docker-gc#manual-usage
# m h dom mon dow command
0 * * * * /usr/bin/gc-job.sh
